# conf/systemd.service
[Unit]
Description=PrivTracker: Private BitTorrent tracker generator
After=network.target
Requires=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
RestartSec=1s
Restart=on-failure
Environment=PORT=__PORT__
AmbientCapabilities=CAP_NET_BIND_SERVICE
WorkingDirectory=__INSTALL_DIR__/
ExecStart=__INSTALL_DIR__/privtracker

[Install]
WantedBy=multi-user.target
